# docker-config/fuseki/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log

[program:executor]
startretries=5
command=/bin/sh /start.sh
stdout_logfile=/var/log/fuseki.log
redirect_stderr=true

[inet_http_server]
port = 9001
username = obeu
password = obeu
