# docker-config/virtuoso/supervisord.conf
[supervisord]
nodaemon=true

[program:virtuoso]
command=/bin/sh /startVirtuoso.sh
startretries=10
startsecs=10