# docker-config/virtuoso/supervisord.conf
[supervisord]
nodaemon=true

[program:virtuoso]
command=/bin/sh /startVirtuoso.sh
startretries=10
redirect_stderr=true
stdout_logfile=/data/virtuoso_supervisor.log
autorestart=true
autostart=true